# worker/macos/org.halide-lang.buildbot.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>org.halide-lang.buildbot</string>
        <key>Program</key>
        <string>${WORKER_SCRIPT}</string>
        <key>RunAtLoad</key>
        <true/>
        <key>EnvironmentVariables</key>
        <dict>
            <key>HALIDE_BB_WORKER_NAME</key>
            <string>${HALIDE_BB_WORKER_NAME}</string>
            <key>PATH</key>
            <string>${LAUNCHD_PATH}</string>
            <key>HL_WEBGPU_NODE_BINDINGS</key>
            <string>${HL_WEBGPU_NODE_BINDINGS}</string>
            <key>HL_WEBGPU_NATIVE_LIB</key>
            <string>${HL_WEBGPU_NATIVE_LIB}</string>
            <key>EMSDK</key>
            <string>${EMSDK}</string>
            <key>HALIDE_NODE_JS_PATH</key>
            <string>${HALIDE_NODE_JS_PATH}</string>
        </dict>
    </dict>
</plist>
